# service/vpnrotate/run
#!/bin/sh

set -e
exec 2>&1

sleep 1

exec vpnrotate --resources=/etc/vpnrotate/resources
